# crontab
0 3 * * * /root/backup-multicraft-servers.sh
0 1 * * * /root/delete-old-backups.sh
