# crontab
0 3 * * * /root/backup-multicraft-servers.sh
